# scrollback.service
[Unit]
Description=Scrollback Keys

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/scrollback -k

[Install]
WantedBy=default.target
